# ansible/roles/plex-backup/templates/plex-backup.timer.j2
[Unit]
Description=Daily Plex Database Backup Timer
Requires=plex-backup.service

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target

